# hws/HW-04/src/sysinfo.service
[Unit]
Description=System Info Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/sysinfo.sh
StandardOutput=append:/var/log/sysinfo.log
StandardError=append:/var/log/sysinfo.log
#Restart=always

[Install]
WantedBy=multi-user.target